# canusb/canusb.rules
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0403", ATTRS{idProduct}=="ffa8", RUN+="/lib/udev/canusb"

